# etc/supervisord.d/screenshot.ini
[program:screenshot]
command=/opt/qnib/phantomjs/bin/screenshot.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart = false
autostart = true
startsecs=1
startretries=0
exitcodes=0
